# packages/server/postinst
#!/bin/sh

set -e

case "$1" in
	configure)
		update-rc.d yeti-management defaults
		mkdir -p /var/log/yeti
		;;  
	abort-upgrade|abort-remove|abort-deconfigure)
		exit 0
		;;
	*)
		echo "postinst called with unknown argument \$1'" >&2
		exit 1
		;;
esac

. /usr/share/debconf/confmodule

invoke-rc.d rsyslog restart

exit 0
